# flash_memory.service
[Unit]
Description=Start service vblade

[Service]
User=root
ExecStart=/home/admin/flash_memory

[Install]
WantedBy=multi-user.target
